# etc/systemd/system/goagent.service.d/override.conf
[Service]
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty12
